# deploy/ubuntu/telegram-risk-net-guard.timer
[Unit]
Description=Run Telegram Risk Control Egress Guard every 30s

[Timer]
OnBootSec=45s
OnUnitActiveSec=30s
AccuracySec=5s
Unit=telegram-risk-net-guard.service

[Install]
WantedBy=timers.target
